# backend/systemd/material-management-backend.service
[Unit]
Description=Material Management Backend Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/material-management/backend
Environment="PATH=/opt/material-management/backend/venv/bin"

# 生产环境请将下方 SECRET_KEY 替换为随机字符串（python3 -c "import secrets; print(secrets.token_hex(32))"）
Environment="SECRET_KEY=CHANGE_ME_TO_A_RANDOM_32_CHAR_STRING"

# CORS_ORIGINS：填入你的实际域名，多个用英文逗号分隔
Environment="CORS_ORIGINS=http://keeponline.cn,http://www.keeponline.cn"

# 使用单进程 + 多线程模式，避免 SQLite 多进程写锁冲突
# -w 1 --threads 4：1 个 worker 进程，4 个线程，适合 SQLite 并发场景
ExecStart=/opt/material-management/backend/venv/bin/gunicorn \
    -w 1 \
    --threads 4 \
    -b 127.0.0.1:5000 \
    --timeout 120 \
    --keep-alive 5 \
    --access-logfile - \
    --error-logfile - \
    app:app

Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
